# src/etc_conf/pathogenicityCalculator.crontab
10	01	*	*	*	ruby ${DIR_TARGET}/bin/config2ConclusionCachePrimer.rb ${DIR_TARGET}/conf/pathoCalc/configColl.json       > ${DIR_TARGET}/var/logs/config2ConclusionCachePrimer.out 2> ${DIR_TARGET}/var/logs/config2ConclusionCachePrimer.err; ruby ${DIR_TARGET}/bin/config2EvidenceCachePrimer.rb ${DIR_TARGET}/conf/pathoCalc/configColl.json       > ${DIR_TARGET}/var/logs/config2EvidenceCachePrimer.out 2> ${DIR_TARGET}/var/logs/config2EvidenceCachePrimer.err
